# services/service-plugin/run
#!/bin/bash

set -e

export SERVICE_NAME=service-plugin
export PORT=4004

export NOTION_CLIENT_ID=$NOTION_CLIENT_ID
export NOTION_CLIENT_SECRET=$NOTION_CLIENT_SECRET

export MONGO_URL=$MONGO_URL_LOCAL_PLUGIN
export GITHUB_APP_ID=222246
export GITHUB_APP_WEBHOOK_SECRET=$GITHUB_APP_WEBHOOK_SECRET
export GITHUB_APP_CLIENT_ID=$GITHUB_APP_CLIENT_ID
export GITHUB_APP_CLIENT_SECRET=$GITHUB_APP_CLIENT_SECRET

export GITHUB_OAUTH_APP_CLIENT_ID=$GITHUB_OAUTH_APP_CLIENT_ID
export GITHUB_OAUTH_APP_CLIENT_SECRET=$GITHUB_OAUTH_APP_CLIENT_SECRET
export GITHUB_APP_PRIVATE_KEY=$GITHUB_APP_PRIVATE_KEY

export GCAL_CLIENT_ID=$GCAL_CLIENT_ID
export GCAL_PROJECT_ID=$GCAL_PROJECT_ID
export GCAL_CLIENT_SECRET=$GCAL_CLIENT_SECRET
export GCAL_JAVASCRIPT_ORIGIN=$GCAL_JAVASCRIPT_ORIGIN

go test ./...
go build -tags=jsoniter -o build/service-plugin
./build/service-plugin
